# systemd/pi-audiocast.service
[Unit]
Description=Raspberry Pi DLNA Audio Renderer
After=network.target sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/dlna-renderer
ExecStart=/usr/bin/python3 -m src.main -c config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=512M
CPUQuota=80%

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
